# crontab
* * * * * cd /app && docker-compose run -rm artisan schedule:run >> /dev/null 2>&1
